# packaging/systemd/chronos-gateway.service
[Unit]
Description=Chronos Gateway (ntp_shm, non-root)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Non-privileged by default: the ntp_shm backend writes a world-writable SysV
# SHM segment, so the gateway needs no root. For the chrony_sock backend (which
# must write chronyd's root-owned socket) use chronos-gateway-chrony.service.
DynamicUser=yes
ExecStart=/usr/local/bin/chronos-gateway --config /etc/chronos/gateway.yaml
Restart=always
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
