# packaging/systemd/chronos-server.service
[Unit]
Description=Chronos Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=chronos
Group=chronos
ExecStart=/usr/local/bin/chronos-server --config /etc/chronos/server.yaml
Restart=always
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
