# packaging/systemd/chronos-gateway.service
[Unit]
Description=Chronos Gateway
After=network-online.target chronyd.service
Wants=network-online.target chronyd.service

[Service]
Type=simple
# Runs as root: chronyd creates the SOCK refclock socket owned by root, so the
# gateway must be root to write samples to it.
User=root
Group=root
ExecStart=/usr/local/bin/chronos-gateway --config /etc/chronos/gateway.yaml
Restart=always
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/run/chrony

[Install]
WantedBy=multi-user.target
